# deploy/flappy.service
# Bare-metal deployment (no Docker) via systemd.
#
#   sudo useradd --system --create-home --home-dir /opt/flappy flappy
#   sudo cp -r . /opt/flappy && cd /opt/flappy
#   sudo -u flappy python3 -m venv .venv
#   sudo -u flappy ./.venv/bin/pip install -r requirements.txt
#   sudo cp deploy/flappy.service /etc/systemd/system/
#   sudo systemctl enable --now flappy
#
# Binds to localhost only; put Nginx (deploy/nginx.conf) in front for TLS.

[Unit]
Description=Flappy Bird QR game
After=network.target

[Service]
Type=simple
User=flappy
WorkingDirectory=/opt/flappy
Environment=PORT=8000
# Uncomment to pin the public URL used in the QR code:
# Environment=PUBLIC_BASE_URL=https://game.example.com
ExecStart=/opt/flappy/.venv/bin/uvicorn main:app \
    --host 127.0.0.1 --port 8000 \
    --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=3

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/flappy

[Install]
WantedBy=multi-user.target
